# server-config/systemd/mc-idle-watchdog.service
[Unit]
Description=Shut down the server after 15 minutes without players
After=minecraft.service

[Service]
Type=oneshot
ExecStart=/opt/minecraft/bin/idle-watchdog.sh
