# server-config/systemd/mc-up-notify.service
[Unit]
Description=Notify Discord when the Minecraft world has finished loading
After=minecraft.service
Requires=minecraft.service

[Service]
Type=oneshot
ExecStart=/opt/minecraft/bin/mc-up-notify.sh

[Install]
WantedBy=multi-user.target
